# rtsp-to-webrtc/rootfs/etc/services.d/rtsp-to-webrtc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the RTSPtoWebRTC service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Generate config.json
bashio::log.info "Generating config.json "
mv /etc/services.d/rtsp-to-webrtc/config.json /app/config.json

bashio::log.info "Starting discovery service"
./discovery &

bashio::log.info "Starting RSTPtoWebRTC"
cd /app
exec /app/rtsp-to-webrtc
