# rtsp-to-webrtc/rootfs/etc/services.d/rtsp-to-webrtc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the RTSPtoWebRTC service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting discovery service"
./discovery &

bashio::log.info "Starting RSTPtoWebRTC"
cd /app
exec /app/rtsp-to-webrtc
